# ansible/roles/llm/templates/buoy-llm.service.j2
[Unit]
Description=Buoy LLM containers (Ollama, Whisper, LLM node)
After=buoy-ros.service docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ buoy_root }}/docker
Environment=ROS_DOMAIN_ID={{ ros_domain_id }} DOCKER_PLATFORM={{ docker_platform }} LLM_MODEL={{ llm_model }} LLM_VISION_MODEL={{ llm_vision_model | default('moondream') }}
ExecStart=/usr/bin/docker compose --profile llm up -d
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
