# ansible/roles/command_center/templates/buoy-command-center.service.j2
[Unit]
Description=Buoy Web Portal
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory={{ buoy_root }}/command_center
Environment=COMMAND_CENTER_PORT={{ command_center_port }}
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server.js
Restart=unless-stopped
# Read DHCP leases (often root-only)
User=root

[Install]
WantedBy=multi-user.target
